# crontab
# START CRON JOB
* * * * * python /app/hello.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
