# crontab
#run python script every minutes
* * * * * python /app/hello.py > /proc/1/fd/1 2>/proc/1/fd/2

