# crontab
# START CRON JOB
* * * * * python3 /app/hello.py >>/tmp/out.log 2>/tmp/err.log
# END CRON JOB