# cron-numpy/crontab
#run python script every minutes
* * * * * /usr/local/bin/python /app/main.py > /proc/1/fd/1 2>/proc/1/fd/2
